# service/snake_dance.service
[Unit]
Description=Snake Dance

[Service]
Type=simple
ExecStart=/home/sina/projects/snake/venv/bin/python -m snake.main
