# service/snake_dance.ini
[program:snake_dance]
command=/home/sina/snake-deploy/snake-env/bin/python -m snake.main
directory=/home/sina/snake-deploy
user=sina
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/snake_dance.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PYTHONPATH="/home/sina/snake-deploy/snake-env/lib/python3.9/site-packages"
stopsignal=TERM
stopwaitsecs=10
